# debian/element-web-installer.cron.d
12 4 * * *	root	chronic /usr/sbin/element-web-installer upgrade
