# debian/grouprise.cron.d
# update search index
*/5 * * * * _grouprise        chronic grouprisectl update_index --remove
